# timers/lazyass.service
[Unit]
Description=LazyASS notification service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lazyass.sh
